# debian/kafka-daemon.postinst
#!/bin/sh

set -e

case "$1" in
    configure)
        update-rc.d kafka defaults || :
        chown -R kafka: /var/log/kafka /var/lib/kafka
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
